# 15-deployment/blue_yellow_app_deployment/deploy/roles/common/templates/supervisor_app.conf.j2
[program:{{ app_name }}]
command={{ venv_dir }}/bin/pserve production.ini
directory={{ sub_app_dir }}
user={{ deploy_user }}
autostart=true
autorestart=true
redirect_stderr=True
